# slack/dashcams-check-tag-failures.timer
[Unit]
Description=Check for tag failures daily

[Timer]
OnCalendar=*-*-* 05:30:00

[Install]
WantedBy=timers.target
