# slack/dashcams-check-tag-failures.service
[Unit]
Description=Check for tag failures

[Service]
Type=oneshot
ExecStart=/usr/local/DASHCAMS/dashcams-check-tag-failures.sh
Group=dashcams
